# recipes-qt/qt6/qtinterfaceframework/if-services.target
[Unit]
Description=Qt Interface Framework services
After=systemd-user-sessions.service
